# docker/demo-stack/instances/retriever/mw/supervisor.conf
[program:start.pl]
command=/opt/app/scripts/start.pl
autostart=true
autorestart=false
