# docker/demo-stack/instances/storer/cron/crontab
# cloud crontab (demo-stack, keep IPC-like flow)
*/8 * * * * /opt/app/scripts/sync_sqlite.pl
*/7 * * * * /opt/app/scripts/test/serverproperty.pl -r
*/8 * * * * /opt/app/scripts/sync_s3.pl

# LTMS download
*/5 * * * * python3 /opt/ltms-client/bin/mediasync.py

# LTMS deploy (mediamgr) + wrapper decides whether to restart containers
5 */6 * * * /opt/ltms-client/bin/mediamgr.pl | /bin/bash /app/check_restart.sh
